# supervisor/supervisord.conf
[supervisord]
logfile = /var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /var/run/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false
#environment = KEY1="value1",KEY2="value2"

[unix_http_server]
file = /var/run/supervisord.sock
chmod = 0777
chown= nobody:nobody
username = user
password = 123

[supervisorctl]
serverurl = unix:///var/run/supervisord.sock
#serverurl =  192.168.1.222:8301
username = user
password = 123
prompt = Hello Supervisor!


#[inet_http_server]
#port = 192.168.1.222:8301
#username = user
#password = 123     

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:test-php]
command=/usr/local/bin/php /tmp/test-php.php
autorstart=true
stdout_logfile=/tmp/test-php.log
